# mudra/deploy/systemd/mudra-api.service
[Unit]
Description=MUDRA API Service
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/mudra
ExecStart=/opt/mudra/.venv/bin/python -m backend.run_api
Restart=always
RestartSec=5
User=mudra
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
